# crontest.cron
15,30,45,59 * * * * echo "xgmtest....." 
